# assets/config/supervisor/2-openbox.conf
[program:Openbox]
command=openbox

priority=2
stderr_logfile=/var/log/openbox.err.log

