# assets/config/supervisor/7-selenium-grid.conf
[program:SeleniumServer]
command=java -jar /opt/selenium-server.jar standalone
autorestart=true

priority=7
stderr_logfile=/var/log/selenium-server.err.log
stdout_logfile=/dev/null
