# rootfs/config/supervisor/6-self-ping.conf
[program:Self-ping]
command=python3 /config/self-ping.py
stderr_logfile=/var/log/selfping.err.log
exitcodes=0
startsecs=0
priority=6
